# files/systemd/restart@.service
[Unit]
Description=Restart service

[Service]
Type=oneshot
ExecStart=/usr/bin/systemctl try-restart minecraft@%i.service
